# server-csharp/build-scripts/dvs-cons/dvs-cons.service
[Unit]
Description=dvs-cons
After=network-online.target
Wants=network-online.target

[Service]
# modify when deploy in prod env
User=dvs
Group=dvs

Type=simple
ExecStart=/usr/local/dotnet-sdk/dotnet /usr/local/sunlight/dvs/dvs-cons/DVS.Cons.Api.dll
WorkingDirectory=/usr/local/sunlight/dvs/dvs-cons

Restart=always
RestartSec=1
StartLimitInterval=0

[Install]
WantedBy=multi-user.target
